# udev_rules/99-maccel.rules
ACTION=="add", SUBSYSTEM=="module", DEVPATH=="/module/maccel", RUN+="maccel_param_ownership_and_resets"
